# scripts/phone-server/units/jcode-pair.service
# Compat wrapper for one release — prefer next-code-pair.service.
# Disable this unit after migrating: systemctl --user disable jcode-pair.service && systemctl --user enable next-code-pair.service
[Unit]
Description=COMPAT: redirects to next-code-pair.service (one-release alias)
After=network.target

[Service]
Type=oneshot
ExecStart=/bin/systemctl --user start next-code-pair.service
RemainAfterExit=yes

[Install]
WantedBy=default.target
